# templates/ai.claude-max-api-proxy.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.claude-max-api-proxy</string>

    <key>Comment</key>
    <string>openclaw-local-bridge-macos — claude-max-api-proxy on :{{PORT}}</string>

    <key>ProgramArguments</key>
    <array>
      <string>{{NODE_BIN}}</string>
      <string>{{PROXY_HOME}}/dist/server/standalone.js</string>
      <string>{{PORT}}</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>WorkingDirectory</key>
    <string>{{HOME}}</string>

    <key>StandardOutPath</key>
    <string>{{HOME}}/.openclaw/logs/claude-max-api-proxy.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/.openclaw/logs/claude-max-api-proxy.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>CLAUDE_CODE_ENTRYPOINT</key>
      <string>cli</string>
      <key>HOME</key>
      <string>{{HOME}}</string>
      <key>NODE_USE_SYSTEM_CA</key>
      <string>1</string>
      <key>NODE_EXTRA_CA_CERTS</key>
      <string>/etc/ssl/cert.pem</string>
      <key>PATH</key>
      <string>{{PATH}}</string>
    </dict>
  </dict>
</plist>
